# context/resources/R/etc/services.d/jupyter-notebook/run
#!/usr/bin/with-contenv bash
jupyter notebook --no-browser --ip 0.0.0.0 --allow-root --port 8888 /home/rstudio